# systemd/system/breezy-renderer.service
# Breezy Box KMS renderer — kiosk system service.
#
# Runs as the unprivileged app user but is bound to tty1 (the active VT on
# seat0), so logind registers it as the seat's session and grants it DRM
# master.  A --user service cannot do this: it lives in the seatless
# user@.service manager scope and drmSetMaster() is denied there.
#
# Confirmed healthy in foreground: the renderer never self-exits; a black screen
# after the overlay briefly appears means it was SIGTERM'd by getty re-spawning
# on tty1 — hence the mask + Restart=always below.
# This unit OWNS tty1.  setup_system.sh MASKS getty@tty1.service so getty can
# never be (re)spawned on that VT — relying on Conflicts= alone is not enough,
# because VT-activation logic can re-pull getty@tty1 and then SIGTERM us as the
# conflict resolution, killing a working renderer (overlay flashes, then black).
# The headless Breezy services (Xvfb, web, VNC, desktop UI) remain --user units
# in breezy.target; only the renderer needs the seat.
[Unit]
Description=Breezy Box KMS renderer (kiosk on tty1)
Documentation=https://github.com/wheaney/breezy-box
After=breezy-gadget.service systemd-user-sessions.service getty@tty1.service
Wants=breezy-gadget.service
# Cap restarts so a stray-signal loop can't hot-spin (see Restart=always below).
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
# __APP_USER__ is substituted by setup_system.sh.
User=__APP_USER__
PAMName=login
# Bind to tty1 so logind seats this session on seat0 (required for DRM master).
# No TTYReset / TTYVHangup: those send hangup signals on the tty that can race
# with VT activation and terminate the renderer.
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
# Groups needed for DRM/GBM, GPU render, and reading input/IMU devices.
SupplementaryGroups=video render input
# The renderer puts its render thread on SCHED_FIFO so head-anchoring stays
# smooth while decode threads load the other cores.  An unprivileged service
# needs CAP_SYS_NICE to set RT scheduling, and LimitRTPRIO raised above the
# priority it requests (20); without these the boost silently no-ops.
AmbientCapabilities=CAP_SYS_NICE
LimitRTPRIO=40
ExecStart=/home/__APP_USER__/.local/bin/displaylink_kms_renderer
# Always restart, even on a clean (SIGTERM) exit, so a stray VT signal can't
# leave the screen black.  The start-limit (in [Unit]) prevents a hot loop.
Restart=always
RestartSec=2
UtmpIdentifier=tty1
UtmpMode=user

[Install]
WantedBy=multi-user.target
