# systemd/system/breezy-gadget.service
[Unit]
Description=Breezy Box OTG RNDIS gadget setup
Documentation=https://github.com/wheaney/breezy-box
# Run after the USB subsystem is up but before user sessions start.
After=systemd-udevd.service
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/lib/breezy-box/breezy_gadget_setup.sh
ExecStop=/bin/sh -c 'echo "" > /sys/kernel/config/usb_gadget/breezy-composite/UDC 2>/dev/null || true'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
