# systemd/user/breezy-x11vnc.service
[Unit]
Description=Breezy Box x11vnc server (Xvfb :1 -> 5900)
Requires=breezy-xvfb.service
After=breezy-xvfb.service
PartOf=breezy.target

[Service]
ExecStart=/usr/bin/x11vnc \
    -display :1 \
    -nopw -shared -forever -loop \
    -defer 10 \
    -xkb
Restart=on-failure
RestartSec=2

[Install]
WantedBy=breezy.target
