# systemd/user/breezy-renderer.service
[Unit]
Description=Breezy Box KMS renderer
Documentation=https://github.com/wheaney/breezy-box
# Soft dep: renderer starts and shows XR_DRIVER_DOWN overlay if xr-driver isn't up.
Wants=xr-driver.service
After=xr-driver.service
PartOf=breezy.target

[Service]
# Wait for breezy-gadget.service (a system-level unit) to finish before starting.
# User units cannot use After= to order against system units, so we poll here.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do systemctl is-active --quiet breezy-gadget.service && exit 0; sleep 1; done; echo "breezy-renderer: timed out waiting for breezy-gadget.service" >&2; exit 1'
ExecStart=%h/.local/bin/displaylink_kms_renderer
Restart=on-failure
RestartSec=2
# CAP_NET_ADMIN is needed to bring up the usb0 interface after the system-level
# gadget setup (breezy-gadget.service) has already configured configfs as root.
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=breezy.target
