# systemd/user/breezy-novnc.service
[Unit]
Description=Breezy Box noVNC proxy (5900 -> 8080)
Requires=breezy-x11vnc.service
After=breezy-x11vnc.service
PartOf=breezy.target

[Service]
ExecStart=/usr/share/novnc/utils/novnc_proxy \
    --vnc localhost:5900 \
    --listen 8080
Restart=on-failure
RestartSec=2

[Install]
WantedBy=breezy.target
