# systemd/user/breezy.target
[Unit]
Description=Breezy Box services
# The KMS renderer is a SYSTEM service (breezy-renderer.service on tty1) so it
# can take DRM master on seat0; it is not part of this user target.  These are
# the headless services that run in the user session.
Wants=breezy-xvfb.service breezy-web.service
After=breezy-xvfb.service breezy-web.service

[Install]
WantedBy=default.target
