# systemd/system/getty@tty1.service.d/autologin.conf
# Autologin the breezy user on tty1 so the user systemd session starts at boot.
# This VT becomes the active VT; the KMS renderer takes DRM master on it.
#
# Replace "breezy" with the actual username if different.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin breezy --noclear %I $TERM
