# systemd/system/getty@tty1.service.d/autologin.conf
# Autologin template — "breezy" is replaced with APP_USER by setup_system.sh.
# This VT becomes the active VT; the KMS renderer takes DRM master on it.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin breezy --noclear %I $TERM
